# config/shared/.config/systemd/user/hypr-update-status.service
[Unit]
Description=Refresh Hypr update status cache
Documentation=man:systemd.service(5)

[Service]
Type=oneshot
ExecStart=%h/.local/bin/hypr-update-status refresh
